# meta-ublox-modules/recipes-kernel/vera-p1-driver-llc-remote/files/vera-p1.rules
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="1fc9", ATTRS{idProduct}=="0102", RUN+="/opt/vera-p1/load.sh"
